# AutoValidation/autoval.cron
# min     hour    day     month   dayweek command
  0       */2     *       *       *       autoval_ExpressCosmics.sh > autoval_ExpressCosmics.log
  30      */2     *       *       *       autoval_ExpressCosmics_merge.sh > autoval_ExpressCosmics_merge.log

